# panel.service
[Unit]
Description=Panel VPS Bot Manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/panel
ExecStart=/usr/bin/node /home/panel/server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=panel-vps
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
